# agents/clawdia/skills/voice-trivia/trivia-voice.service
[Unit]
Description=Trivia Voice Server
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/clawdia-hertz-openclaw/agents/clawdia/skills/voice-trivia/server
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=5
EnvironmentFile=%h/.openclaw/.env

[Install]
WantedBy=default.target
