# agents/clawdia/skills/trivia_voice/trivia-voice.service
[Unit]
Description=Trivia Voice Server
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/openclaw/agents/main/skills/voice-trivia/server
# Override TRIVIA_VOICE_WORKDIR in ~/.openclaw/.env if your repo is elsewhere
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=5
EnvironmentFile=%h/.openclaw/.env

[Install]
WantedBy=default.target
